# scripts/com.pace.session-insight.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.amit.pace-session-insight</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>REPLACE_WITH_HOME/.claude/session-insight.py</string>
        <string>--batch</string><string>6</string>
        <string>--days</string><string>3</string>
    </array>
    <key>StartInterval</key><integer>1200</integer>
    <key>RunAtLoad</key><true/>
    <key>StandardOutPath</key><string>/tmp/pace-session-insight.out</string>
    <key>StandardErrorPath</key><string>/tmp/pace-session-insight.err</string>
    <key>ProcessType</key><string>Background</string>
</dict>
</plist>
