# scripts/com.pace.codex-usage-poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pace.codex-usage-poll</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <!-- Edit this path to where you cloned Pace -->
        <string>REPLACE_WITH_HOME/pace/scripts/codex-usage-poll.py</string>
    </array>
    <key>StartInterval</key>
    <integer>30</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/pace-codex-usage-poll.err</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
